# crontab/crontab-destination
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.WY2lO6/crontab installed on Mon Feb 24 10:04:16 2020)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
# Edit this file to introduce tasks to be run by cron.
# 
# Each task to run has to be defined through a single line
# indicating with different fields when the task will be run
# and what command to run for the task
# 
# To define the time you can provide concrete values for
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').
# 
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
# 
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
# 
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
# 
# For more information see the manual pages of crontab(5) and cron(8)
# 
# m h  dom mon dow   command
#@reboot /bin/sh -c 'sleep 30 && docker start jellyfin 2>&1'
@reboot /usr/bin/rclone mount gcrypt-usmba:p0ds0smb /mnt/gcrypt-usmba --daemon --vfs-cache-mode off --dir-cache-time 30m --read-only
@reboot /usr/bin/rclone mount http-home: /mnt/http-home --daemon --vfs-cache-mode off --vfs-read-chunk-size 10M --vfs-read-chunk-size-limit 100M --dir-cache-time 30m --read-only
#@reboot /bin/sh -c 'sleep 30 && screen -dmS docker_logs docker logs --follow --tail 100 jellyfin'
45 */12 * * * sync && echo 1 | tee /proc/sys/vm/drop_caches 2>&1 | logger -t 'drop_caches'
@reboot sleep 30 && /root/scripts/gen-media-lst.sh
22 */2 * * * /root/scripts/gen-media-lst.sh
33 */12 * * * ls -fR /usr/local/jellyfin/media/
