# crontab/crontab-source
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.KquXkBXmTh installed on Thu Feb 27 01:01:00 2020)
# (Cron version -- $FreeBSD: releng/11.2/usr.sbin/cron/crontab/crontab.c 321241 2017-07-19 20:22:16Z ngie $)
#m  h  dom  mon  dow
@reboot /root/scripts/rclone-serve-http.sh
@reboot /root/scripts/rclone-sync-p0ds0smb.sh
40 */4 * * * /root/scripts/rclone-sync-video.sh
15 1 * * * /root/scripts/rclone-sync-p0ds0smb.sh
